# service/postfix/run
#!/bin/bash

exec /usr/lib/postfix/master

